# timers/backintime.service
[Unit]
Description=daily run of backintime
Wants=backintime.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/backintime backup


