# timers/flatpak.service
[Unit]
Description=daily run of flatpak update
Wants=flatpak.timer

[Service]
Type=oneshot
ExecStart=flatpak update -y


